# gameserveros/systemd/gameserveros-audit-forward.timer
[Unit]
Description=Hourly audit digest for agent pickup (future heartbeat field)

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
